# roles/borg_create/templates/borg-create.timer
[Unit]
Description=Borg Backup Timer

[Timer]
OnCalendar={{ borg_create_schedule }}

[Install]
WantedBy=timers.target
