# roles/borg_maintenance/templates/borg-compact.timer
[Unit]
Description=Borg Backup Compact Timer

[Timer]
OnCalendar={{ borg_compact_schedule }}

[Install]
WantedBy=timers.target
