# roles/borg_maintenance/templates/borg-check.timer
[Unit]
Description=Borg Backup Check Timer

[Timer]
OnCalendar={{ borg_check_schedule }}

[Install]
WantedBy=timers.target
