# roles/borg_maintenance/templates/borg-prune.timer
[Unit]
Description=Borg Backup Prune Timer

[Timer]
OnCalendar={{ borg_prune_schedule }}

[Install]
WantedBy=timers.target
